# crontab/crontab_command
# /home/k150030/working 디렉토리를 벤자민 빌드 서버의 backup 폴더로 백업
# 매주 금요일 12:00에 백업
# 0분 12시 모든 일 모든월 금요일 실행프로그램
0 12 * * 5 /home/k150030/log/backup/backup.sh

